# deploy/mac/com.elvira.mac-memory-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

  <key>Label</key>
  <string>com.elvira.mac-memory-worker</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>-m</string>
    <string>mac_memory_guard.worker_agent</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/elvira/scripts</string>

  <key>StartInterval</key>
  <integer>20</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/elvira/logs/mac_worker.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/elvira/logs/mac_worker.err.log</string>

</dict>
</plist>
